# static/files/agent.mpkg/Contents/Packages/internal.pkg/Contents/Resources/meshagent_osx64_LaunchAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>meshagent</string>
    <key>ProgramArguments</key>
    <array>
		  <string>/usr/local/mesh_services/meshagent/meshagent_osx64</string>
		  <string>-kvm1</string>
    </array>
	  <key>LimitLoadToSessionType</key>
	  <array>
		  <string>LoginWindow</string>
	  </array>
    <key>WorkingDirectory</key>
    <string>/usr/local/mesh_services/meshagent/</string>
	  <key>RunAtLoad</key>
	  <true/>
    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
